# workshift-calendar/rootfs/etc/services.d/workshift/finish
#!/usr/bin/with-contenv bashio

EXIT_CODE=$1
SIGNAL_NUM=$2

bashio::log.info "WorkShift Calendar stopped - exit code: ${EXIT_CODE}, signal: ${SIGNAL_NUM}"

if [ "$EXIT_CODE" -eq 256 ] && [ "$SIGNAL_NUM" -eq 15 ]; then
    bashio::log.info "Normal shutdown via SIGTERM"
    exit 0
fi

if [ "$EXIT_CODE" -ne 0 ]; then
    bashio::log.error "Abnormal exit detected - halting container"
    /run/s6/basedir/bin/halt
fi
