# workshift-calendar/rootfs/etc/services.d/workshift/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WorkShift Calendar v1.2.0..."

ADMIN_MODE=$(bashio::config 'admin_mode')
export ADMIN_MODE

if bashio::config.has_value 'admin_mode'; then
    bashio::log.info "Admin mode: $ADMIN_MODE"
fi

export SESSION_SECRET=$(cat /proc/sys/kernel/random/uuid)

export DATABASE_URL="sqlite:////data/workshift.db"
bashio::log.info "Using SQLite database at /data/workshift.db"

cd /app
bashio::log.info "Starting gunicorn on port 5000..."
exec python3 -m gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    app:app
