# workshift-calendar/rootfs/etc/services.d/workshift/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WorkShift Calendar v1.4.0..."

PORT=$(bashio::config 'port')
export PORT
bashio::log.info "Server port: ${PORT}"

LOG_LEVEL=$(bashio::config 'log_level')
export LOG_LEVEL
bashio::log.info "Log level: ${LOG_LEVEL}"

export SESSION_SECRET=$(cat /proc/sys/kernel/random/uuid)

export DATABASE_URL="sqlite:////data/workshift.db"
bashio::log.info "Using SQLite database at /data/workshift.db"

cd /app

bashio::log.info "Starting gunicorn on port ${PORT}..."
exec python3 -m gunicorn \
    --bind "0.0.0.0:${PORT}" \
    --workers 1 \
    --threads 2 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    --log-level "${LOG_LEVEL}" \
    app:app
